# backup.cron
#!/bin/sh
/usr/share/opi-backup/opi-backup.sh -p || exit 0