# backup.cron
#!/bin/sh
/usr/share/opi-backup/opi_backup.sh -p || exit 0
